# scripts/_cloudflared-tunnel.service
[Unit]
Description=Cloudflare Quick Tunnel for music-bot TMA
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel --url http://localhost:8080 --no-autoupdate
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
